# launchagents/com.local.worldpackers-tracker.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN"
  "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>com.local.worldpackers-tracker</string>
  <key>ProgramArguments</key>
  <array>
    <string>/bin/sh</string>
    <string>/Users/soob/Developer/worldpackers-hours-tracker/bin/start.sh</string>
  </array>
  <key>RunAtLoad</key>
  <true/>
  <key>KeepAlive</key>
  <false/>
  <key>StandardOutPath</key>
  <string>/Users/soob/Library/Logs/worldpackers-launchagent.log</string>
  <key>StandardErrorPath</key>
  <string>/Users/soob/Library/Logs/worldpackers-launchagent.log</string>
</dict>
</plist>
